# deploy/marketing-tool.service
[Unit]
Description=Halogen Marketing Automation Platform
After=network.target postgresql.service

[Service]
Type=simple
User=marketing
WorkingDirectory=/opt/marketing-tool
EnvironmentFile=/opt/marketing-tool/.env
ExecStart=/opt/marketing-tool/.venv/bin/uvicorn portal.api.app:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
